# ansible/roles/c2-admin-node/templates/slurm-db-backup.service.j2
[Unit]
Description=Slurm db backup service
Wants=slurm-db-backup.timer

[Service]
Type=oneshot
User=root
Group=root
ExecStart=-/admin/bin/slurm-db-backup
