# ansible/roles/c2-admin-node/templates/slurm-db-backup.timer.j2
[Unit]
Description=Slurm db backup timer
Requires=slurm-db-backup.service

[Timer]
Unit=slurm-db-backup.service
# daily at midnight
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
